# cookbooks/tenkai-suru-api/files/default/tenkai_suru_api.conf
description "Run tenkai-suru-api"

start on runlevel startup
stop on runlevel shutdown

respawn

script
  cd /var/tenkai_suru_api/
  lein with-profile dev trampoline run 8080
end script
